# deploy/terminator.service
[Unit]
Description=Terminator Trading Server
After=network-online.target tailscaled.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=fw
WorkingDirectory=/home/fw/terminator_prod/server
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/bash /home/fw/terminator_prod/deploy/fetch_secrets.sh
ExecStart=/home/fw/terminator_prod/.venv/bin/python main.py
Restart=always
RestartSec=10
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
