# deploy/terminator-downloader.service
[Unit]
Description=Terminator Data Downloader
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=fw
WorkingDirectory=/home/fw/terminator_prod/server/downloader
ExecStartPre=/bin/bash /home/fw/terminator_prod/deploy/fetch_secrets.sh
ExecStart=/home/fw/terminator_prod/.venv/bin/python downloader.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
